# .config/systemd/user/rclone-gdrive.service
[Unit]
Description=rclone mount for Google Drive
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStartPre=/usr/bin/mkdir -p %h/gdrive
ExecStart=/home/linuxbrew/.linuxbrew/bin/rclone mount google-drive: %h/gdrive --vfs-cache-mode full
ExecStop=/bin/fusermount -u %h/gdrive
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
